# client-side/conf/app_supervisor.conf
[program:flaskapp]
command = uwsgi --ini /docker/conf/app_uwsgi.ini

[program:nginx]
command = service nginx restart
